# ttylinux-build/pkg-cfg/iptables-1.4.12/files
etc/firewall.conf
etc/rc.d/init.d/firewall
etc/rc.d/rc.startup/01.firewall
etc/sysconfig/firewall
lib/xtables
sbin/iptables
sbin/xtables-multi
