# ttylinux-build/pkg-cfg/ttylinux-basefs-1.0/files
bin
boot
dev
dev/console
dev/null
etc
etc/cron.d
etc/cron.daily
etc/cron.hourly
etc/cron.monthly
etc/cron.weekly
home
home/user
lib
mnt
mnt/a:
mnt/c:
mnt/cdrom
mnt/flash
mnt/floppy
mnt/usbdisk
proc
root
sbin
sys
tmp
usr
usr/bin
usr/etc
usr/lib
usr/libexec
usr/sbin
usr/share
usr/share/ttylinux
usr/tmp
usr/var
var
var/crash
var/lib
var/lib/misc
var/lock
var/lock/subsys
var/log
var/nobody
var/run
var/spool
var/spool/cron
var/spool/cron/crontabs
var/tmp
