# ttylinux-build/pkg-cfg/thttpd-2.25b/files
etc/rc.d/init.d/thttpd
etc/thttpd.conf
home/thttpd
usr/bin/makeweb
usr/sbin/thttpd
